# tools/inventree-part-import/com.openskyhawk.inventree-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd agent: weekly InvenTree part refresh (Mon 03:00).
     Install: edit the ProgramArguments path to this repo, then
       cp this file ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.openskyhawk.inventree-weekly.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openskyhawk.inventree-weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/ABSOLUTE/PATH/TO/tools/inventree-part-import/weekly-update.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>3</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
